# rootfs/etc/services.d/pw-feeder/run
#!/usr/bin/env bash
#shellcheck shell=bash

/usr/local/sbin/pw-feeder
sleep 10
