# rootfs/etc/services.d/mlat_tunnel_r/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

# Start ACARS tunnel if required
if [[ -n "${ACARS_HOST}" ]]; then
    # shellcheck disable=SC2016
    stunnel /etc/stunnel/stunnel_mlat_results.conf \
        2>&1 | stdbuf -o0 awk '{print "[mlat_tunnel_r] " strftime("%Y/%m/%d %H:%M:%S", systime()) " " $0}'
else
    sleep 86500
fi